# smokeping.ini
[program:smokeping]
command=/usr/bin/smokeping --nodaemon
autostart=true
autorestart=true
startsecs=60
priority=1
stopasgroup=true
killasgroup=true
